# lib/systemd/system/certs-update.timer
[Unit]
Description=Scheduled certificate renewal

[Timer]
OnCalendar=daily
RandomizedDelaySec=2h
Persistent=true

[Install]
WantedBy=timers.target
