# lib/systemd/system/nextcloudcron.timer
[Unit]
Description=Background tasks for Nextcloud

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
