# lib/systemd/system/postgres-backup@.timer
[Unit]
Description=Creates a backup of postgres database %i

[Timer]
OnBootSec=15min
OnUnitActiveSec=12h

[Install]
WantedBy=timers.target
