# lib/systemd/system/unit-status-mailer@.service
[Unit]
Description=Unit Status Mailer Service
Wants=network.target
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/unit-status-mail %i "Hostname: %H" "Machine ID: %m" "Boot ID: %b"
