# ansible/roles/restic-backup/templates/directory-backup-check@.timer.j2
[Unit]
Description=Checks snapshots of repository %i

[Timer]
OnCalendar=06/12:05

[Install]
WantedBy=timers.target
