# ansible/roles/restic-backup/templates/directory-backup@.timer.j2
[Unit]
Description=Creates a backup of repository %i

[Timer]
OnBootSec=15min
OnUnitActiveSec=12h

[Install]
WantedBy=timers.target
