# ansible/roles/web1090/templates/web1090cli-daily.timer
[Unit]
Description=Stores ADS-B data in Postgres.

[Timer]
OnBootSec=15min
OnUnitActiveSec=24h

[Install]
WantedBy=timers.target
